# root/etc/services.d/homeassistant/run
#!/usr/bin/with-contenv bash
exec \
	s6-setuidgid abc \
		python3 -m homeassistant --skip-pip --config /config --log-rotate-days 7 >/dev/stdout 2>/dev/stderr
